# templates/agent-loop/.claude/services/claude-agent.timer
[Unit]
Description=Run Claude Agent Heartbeat every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
